# src/dot_config/systemd/user/9router.service
[Unit]
Description=9router
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/.9router
ExecStart=%h/.local/share/mise/installs/npm-9router/latest/bin/9router
Environment=PATH=%h/.local/share/mise/installs/node/lts/bin:%h/.local/share/mise/installs/npm-9router/latest/bin:/usr/bin:/bin:/usr/sbin:/sbin
Restart=on-failure
RestartSec=5
TimeoutStopSec=60
KillMode=control-group
StandardOutput=append:%h/.9router/logs/stdout_9router.log
StandardError=append:%h/.9router/logs/stderr_9router.log

[Install]
WantedBy=default.target
